# reservation-push-app/reservation-app-project/gunicorn.service
[Unit]
Description = reservation-app-service
After = network.target

[Service]
User=root
Group=root
WorkingDirectory